# rootfs/usr/lib/opkg/info/polipo.list
/usr/sbin/polipo
/etc/init.d/polipo
/etc/config/polipo
/etc/polipo/config
