# rootfs/usr/lib/opkg/info/sshtunnel.list
/etc/config/sshtunnel
/etc/init.d/sshtunnel
